# irssi-libnotify.service
[Unit]
Description=irssi-libnotify-git

[Service]
ExecStart=/usr/bin/notify-listener.py

[Install]
WantedBy=default.target
